# SystemD/backend.service
[Unit]
Description=Temperature and Humidity Sensor Backend
After=network.target

[Service]
ExecStart=/home/%u/AHT21/environment/bin/python3 /home/%u/AHT21/backend.py
WorkingDirectory=/home/%u/AHT21
Restart=always
User=%u
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
